# profile.d/testd.sh
function testd {
    guid="$(uuidgen)"
    mkdir "/tmp/${guid}"
    cd "/tmp/${guid}"
}

